# roles/unifi/templates/unifi-controller.service.j2
# container-unifi-controller.service
# autogenerated by Podman 2.2.1
# Tue Jan 19 12:13:13 UTC 2021
[Unit]
Description=Podman container-unifi-controller.service
Documentation=man:podman-generate-systemd(1)
Wants=network.target
After=network-online.target
[Service]
Environment=PODMAN_SYSTEMD_UNIT=%n
Restart=on-failure
ExecStartPre=/bin/rm -f %t/container-unifi-controller.pid %t/container-unifi-controller.ctr-id
ExecStart=/usr/bin/podman run --conmon-pidfile %t/container-unifi-controller.pid --cidfile %t/container-unifi-controller.ctr-id --cgroups=no-conmon -d --replace --name unifi-controller -e PUID=1000 -e PGID=1000 -e MEM_LIMIT=1024 -v {{ unifi_user_info.home }}/{{ unifi_unifi_config_dir }}:/config:Z --label io.containers.autoupdate=image -p 3478:3478/udp -p 10001:10001/udp -p 8082:8080 -p 8443:8443 -p 8843:8843 -p 8880:8880 -p 6789:6789 -p 5514:5514 ghcr.io/linuxserver/unifi-controller:{{ unifi_controller_version }}
ExecStop=/usr/bin/podman stop --ignore --cidfile %t/container-unifi-controller.ctr-id -t 10
ExecStopPost=/usr/bin/podman rm --ignore -f --cidfile %t/container-unifi-controller.ctr-id
PIDFile=%t/container-unifi-controller.pid
KillMode=none
Type=forking
TimeoutSec=900
[Install]
WantedBy=multi-user.target default.target
